# example.crontab
DISPLAY=:0
XDG_RUNTIME_DIR=/run/user/10001
00,30 05,06,07,08 * * * /usr/local/bin/nh-weather-closings.pl
